# deploy_tools/gunicorn-upstart.template.conf
description "Gunivorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid strngluv
chdir /home/strngluv/sites/SITENAME/source

exec ../vituralenv/bin/gunicorn \
	--bind unix:/tmp/SITENAME.socket \
		superlists.wsgi:application